# bot.service
[Unit]
Description=Telegram Flood Bot
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/bot
ExecStart=/home/YOUR_USER/bot/venv/bin/python main.py
Restart=always
RestartSec=5
EnvironmentFile=/home/YOUR_USER/bot/.env

[Install]
WantedBy=multi-user.target
